# deploy/opswatch.service
[Unit]
Description=OpsWatch self-hosted ops stack (scheduler, monitors, alerts)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=/var/lib/opswatch
Environment=PYTHONPATH=__APP_DIR__
Environment=OPSWATCH_STORE_PATH=/var/lib/opswatch/opswatch.db
Environment=OPSWATCH_STATE_DIR=/var/lib/opswatch
# Secrets live here, never in the repo or the config file. Uncomment the ones
# you use. Each maps to a channel or feature declared in the config:
# Environment=OPSWATCH_WEBHOOK_URL=https://hooks.example.com/your/endpoint
# Environment=OPSWATCH_SLACK_URL=https://hooks.slack.com/services/XXX/YYY/ZZZ
# Environment=OPSWATCH_TELEGRAM_TOKEN=123456:abcdef
# Environment=OPSWATCH_TELEGRAM_CHAT_ID=987654321
# Environment=OPSWATCH_SMTP_USER=ops@example.com
# Environment=OPSWATCH_SMTP_PASSWORD=app-password-here
# Dashboard login (optional, set "auth.enabled": true in the config too):
# Environment=OPSWATCH_DASHBOARD_PASSWORD=choose-a-strong-one
# Token external jobs use to post heartbeats and failures to /api/ingest:
# Environment=OPSWATCH_INGEST_TOKEN=long-random-string
ExecStart=/usr/bin/python3 -m opswatch --config __CONFIG__
Restart=on-failure
RestartSec=5

# Hardening: the service needs almost nothing on the box.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/opswatch

[Install]
WantedBy=multi-user.target
